# crontab
*/30 * * * * python /main.py >> /result
